# build-debian/avertem/lib/systemd/system/avertemd.service
[Unit]
Description=avetemd
[Service]
ExecStart=/opt/avertem/bin/avertemd.sh
PrivateTmp=true
