# build-debian/avertem/etc/systemd/system/avertemd.service
[Unit]
Description=Avertem

[Service]
ExecStart=/opt/avertem/bin/avertemd.sh

[Install]
WantedBy=multi-user.target
